# rss-daily-digest/systemd/rss-digest.service
[Unit]
Description=RSS Daily Digest
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/home/xkong/rss-daily-digest
ExecStart=/home/xkong/rss-daily-digest/run.sh
