# rss-daily-digest/systemd/rss-digest.timer
[Unit]
Description=Run RSS Digest daily at 09:00 CST (01:00 UTC)

[Timer]
OnCalendar=*-*-* 01:00:00
Persistent=true

[Install]
WantedBy=timers.target
